# bot.conf
[program:biplexsystems]
directory=/root/bot
command=env/bin/python3 app.py
autorestart=true
autostart=true
